# build/minecraft-server.conf
[program:minecraft-server]
user = nobody
command = /home/nobody/start.sh
umask = 000
