# build/minecraft-server.conf
[program:start-script]
autorestart = false
startsecs = 0
user = nobody
command = /home/nobody/start.sh
umask = 000

[program:purge-script]
autorestart = false
startsecs = 0
user = nobody
command = /home/nobody/purge.sh
umask = 000

[program:backup-script]
autorestart = false
startsecs = 0
user = nobody
command = /home/nobody/backup.sh
umask = 000

[program:shutdown-script]
autorestart = false
startsecs = 0
user = root
command = /usr/local/bin/shutdown.sh '^java'
umask = 000